# systemd/jobsonar-analyst.service
[Unit]
Description=JobSonar Analyst Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/JobSonar
ExecStart=/usr/bin/python3 -u main_analyst.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
