# web/jobsonar-cs.service
[Unit]
Description=JobSonar Customer Service Bot (Head of Customer Service)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/JobSonar/web
ExecStart=/usr/bin/python3 -u main_cs.py
Restart=always
RestartSec=10
StandardOutput=append:/root/JobSonar/web/log_cs_systemd.txt
StandardError=append:/root/JobSonar/web/log_cs_systemd.txt

[Install]
WantedBy=multi-user.target
